# sensor/app/crontab
* * * * * /usr/local/bin/python2.7 /app/sensor.py
# An empty line is required at the end of this file for a valid cron file.
